# deployment/banix_webapp.ini
[program:banix_webapp]
command=/opt/banix/banix_server/.venv/bin/gunicorn --bind 0.0.0.0:7701 wsgi:app
directory=/opt/banix/banix_server/web_service
stdout_logfile=/home/administrator/banix/webapp.log
