# deployment/banix_media_proxy.ini
[program:banix_media_proxy]
command=/opt/banix/banix_server/.venv/bin/python /opt/banix/banix_server/src/generate_media_proxy.py
directory=/opt/banix/banix_server/src
stdout_logfile=/home/administrator/banix/proxy.log
